# examples/com.syscleaner.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.syscleaner.weekly</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd "$HOME" &amp;&amp; uv run syscleaner scan --all --save-results "$HOME/.local/share/syscleaner/weekly-scan.json" &amp;&amp; uv run syscleaner report --from-scan "$HOME/.local/share/syscleaner/weekly-scan.json" -f markdown -o "$HOME/.local/share/syscleaner/weekly-report.md"</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/syscleaner-weekly.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/syscleaner-weekly.err</string>
</dict>
</plist>
